# com.adishmitanka.podcastmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adishmitanka.podcastmonitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/adishmitanka/Local Documents/Code/Gemini Transcribe/monitor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/adishmitanka/Local Documents/Code/Gemini Transcribe</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/adishmitanka/Local Documents/Code/Gemini Transcribe/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adishmitanka/Local Documents/Code/Gemini Transcribe/monitor.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>GOOGLE_API_KEY</key>
        <string>AIzaSyCljZ5ELClQfuLSX62ux_myAZxPgL76Cvk</string>
        <key>WATCH_DIRECTORY</key>
        <string>/Users/adishmitanka/Content/Podcasts</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
